# docker/php/rootfilesystem/etc/supervisor/service.d/laravel_queue.conf
[supervisord]
user = root

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www
command=php artisan queue:work --max-time=3600
user = root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stopwaitsecs=3600